# 4_Others/hackfest2021/Internet of Torments (@arii & @Miu.)/_spookyIoTP1.extracted/squashfs-root-0/usr/lib/opkg/info/hostname-tool.list
/etc/init.d/updateHostname
/usr/bin/updateHostname
