# 4_Others/hackfest2021/Internet of Torments (@arii & @Miu.)/_spookyIoTP1.extracted/squashfs-root-0/usr/lib/opkg/info/onion-helper.list
/usr/sbin/onion-helper
/etc/init.d/onion-helper
