# 4_Others/hackfest2021/Internet of Torments (@arii & @Miu.)/_spookyIoTP1.extracted/squashfs-root-0/usr/lib/opkg/info/shellinabox.list
/usr/lib/shellinabox/style.css
/etc/init.d/shellinabox
/usr/sbin/shellinaboxd
